# build/prowlarr.conf
[program:start-script]
user = nobody
command = /home/nobody/start.sh
umask = 000


